# picoclaw/rootfs/etc/s6-overlay/s6-rc.d/picoclaw-launcher/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -e

export HOME="/share/picoclaw"
export PICOCLAW_HOME="${HOME}/.picoclaw"
mkdir -p "${PICOCLAW_HOME}"

# ── Timezone ──────────────────────────────────────────────────────────────────
TIMEZONE=$(bashio::config 'timezone')
if [ -f "/usr/share/zoneinfo/${TIMEZONE}" ]; then
    ln -sf "/usr/share/zoneinfo/${TIMEZONE}" /etc/localtime
    export TZ="${TIMEZONE}"
fi

# ── Home Assistant MCP ────────────────────────────────────────────────────────
HA_MCP_ENABLED=$(bashio::config 'ha_mcp_enabled')
CONFIG_FILE="${PICOCLAW_HOME}/config.json"

if bashio::var.true "${HA_MCP_ENABLED}"; then
    HA_TOKEN="${SUPERVISOR_TOKEN}"
    MCP_PATCH=$(jq -n --arg token "${HA_TOKEN}" '{
        tools: {
            mcp: {
                enabled: true,
                discovery: { enabled: false },
                servers: {
                    "home-assistant": {
                        enabled: true,
                        type: "http",
                        url: "http://supervisor/core/api/mcp",
                        deferred: false,
                        headers: {
                            Authorization: ("Bearer " + $token)
                        }
                    }
                }
            }
        }
    }')

    if [ -f "${CONFIG_FILE}" ]; then
        jq --argjson mcp "${MCP_PATCH}" '. * $mcp' "${CONFIG_FILE}" > "${CONFIG_FILE}.tmp" \
            && mv "${CONFIG_FILE}.tmp" "${CONFIG_FILE}"
        bashio::log.info "Home Assistant MCP configured."
    else
        echo "${MCP_PATCH}" > "${CONFIG_FILE}"
        bashio::log.info "Home Assistant MCP config stub written (launcher will merge on first start)."
    fi
fi

# ── Browser addon auto-detect ─────────────────────────────────────────────────
BROWSER_CDP_PORT=$(bashio::config 'browser_cdp_port')
BROWSER_CDP_URL="ws://localhost:${BROWSER_CDP_PORT}"
if curl -sf "http://localhost:${BROWSER_CDP_PORT}/json/version" > /dev/null 2>&1; then
    bashio::log.info "Browser addon detected — configuring agent-browser CDP at ${BROWSER_CDP_URL}"
    mkdir -p "${HOME}/.agent-browser"
    printf '{"cdp":"%s"}\n' "${BROWSER_CDP_URL}" > "${HOME}/.agent-browser/config.json"
else
    bashio::log.info "No browser addon detected — browser tools unavailable. Install browserless_chrome addon to enable."
fi

# ── Start launcher ────────────────────────────────────────────────────────────
bashio::log.info "Starting PicoClaw launcher (web UI: port 18800)..."
exec /usr/local/bin/picoclaw-launcher -public -no-browser 2>&1
