# zeroclaw/rootfs/etc/s6-overlay/s6-rc.d/zeroclaw/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "ZeroClaw daemon stopped."
