# zeroclaw/rootfs/etc/s6-overlay/s6-rc.d/zeroclaw/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -e

export HOME="/data/zeroclaw"
mkdir -p "${HOME}/workspace"

# If no config exists yet, run onboard with the provided credentials
if [ ! -f "${HOME}/.zeroclaw/config.toml" ]; then
    bashio::log.info "No config found, running initial onboard..."
    PROVIDER=$(bashio::config 'provider')
    API_KEY=$(bashio::config 'api_key')
    /usr/local/bin/zeroclaw onboard \
        --provider "${PROVIDER}" \
        --api-key "${API_KEY}" \
        --memory markdown \
        --force
    bashio::log.info "Onboard complete. Use the Web UI terminal for further configuration."
fi

bashio::log.info "Starting ZeroClaw daemon..."
exec /usr/local/bin/zeroclaw daemon
