# zeroclaw/rootfs/etc/s6-overlay/s6-rc.d/zeroclaw/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -e

# ── Paths ────────────────────────────────────────────────────────────────────
export HOME="/share/zeroclaw"
export SHELL="/bin/bash"

ZEROCLAW_DATA="/share/zeroclaw/.zeroclaw"
CONFIG_FILE="${ZEROCLAW_DATA}/config.toml"
mkdir -p "${ZEROCLAW_DATA}" "${HOME}/workspace"

# Redirect ZeroClaw's default internal data dir to the persistent share volume.
# v0.6.8+ uses /usr/local/var/zeroclaw/ instead of ~/.zeroclaw/ — symlink it.
mkdir -p /usr/local/var
ln -sfn "${ZEROCLAW_DATA}" /usr/local/var/zeroclaw

# ── Inject provider API keys from .env file into container environment ────────
# Create /share/zeroclaw/.zeroclaw/.env with KEY=value pairs to expose
# additional provider keys (e.g. OPENAI_API_KEY, GROQ_API_KEY) to the daemon.
ENV_FILE="${ZEROCLAW_DATA}/.env"
if [ -f "${ENV_FILE}" ]; then
    while IFS='=' read -r key value || [ -n "$key" ]; do
        [ -z "$key" ] && continue
        case "$key" in \#*) continue ;; esac
        export "${key}=${value}"
        printf '%s' "${value}" > "/var/run/s6/container_environment/${key}" 2>/dev/null || true
    done < "${ENV_FILE}"
    bashio::log.info "Provider env vars loaded from .env"
fi

# ── First-run guard ───────────────────────────────────────────────────────────
if [ ! -f "${CONFIG_FILE}" ]; then
    bashio::log.fatal "No configuration found. Open the web terminal and run: zeroclaw onboard"
    exit 1
fi

# ── Read options ──────────────────────────────────────────────────────────────
TIMEZONE=$(bashio::config 'timezone')

# ── Timezone ──────────────────────────────────────────────────────────────────
if [ -f "/usr/share/zoneinfo/${TIMEZONE}" ]; then
    ln -sf "/usr/share/zoneinfo/${TIMEZONE}" /etc/localtime
    echo "${TIMEZONE}" > /etc/timezone
    bashio::log.info "Timezone set to ${TIMEZONE}"
else
    bashio::log.warning "Unknown timezone '${TIMEZONE}', falling back to UTC."
fi

# ── Remove only the managed 'home-assistant' MCP server block ────────────────
# Uses Node.js (guaranteed available) to strip the [[mcp.servers]] block where
# name = "home-assistant", leaving all other [[mcp.servers]] blocks intact.
cat > /tmp/mcp_strip_ha.js <<'JSEOF'
const fs = require('fs');
const f = process.argv[2];
const lines = fs.readFileSync(f, 'utf8').split('\n');
const out = [];
let buf = [], inServer = false, isHA = false;
for (const line of lines) {
  if (line === '[[mcp.servers]]') {
    if (inServer && !isHA) out.push(...buf);
    buf = [line]; inServer = true; isHA = false;
  } else if (inServer) {
    if (/^\s*name\s*=\s*"home-assistant"/.test(line)) isHA = true;
    if (line.startsWith('[') && line !== '[[mcp.servers]]' && !line.startsWith('[mcp.servers.')) {
      if (!isHA) out.push(...buf);
      buf = []; inServer = false; isHA = false; out.push(line);
    } else { buf.push(line); }
  } else { out.push(line); }
}
if (inServer && !isHA) out.push(...buf);
fs.writeFileSync(f, out.join('\n'));
JSEOF
node /tmp/mcp_strip_ha.js "${CONFIG_FILE}"

# ── Home Assistant MCP integration ────────────────────────────────────────────
HA_MCP_ENABLED=$(bashio::config 'ha_mcp_enabled')
if bashio::var.true "${HA_MCP_ENABLED}"; then
    if grep -q '^\[mcp\]' "${CONFIG_FILE}" 2>/dev/null; then
        sed -i '/^\[mcp\]/,/^\[/{s/^enabled = false$/enabled = true/}' "${CONFIG_FILE}"
        sed -i '/^\[mcp\]/,/^\[/{s/^deferred_loading = true$/deferred_loading = false/}' "${CONFIG_FILE}"
        sed -i '/^servers = \[\]$/d' "${CONFIG_FILE}"
        bashio::log.info "MCP section updated: enabled = true, deferred_loading = false."
    else
        cat >> "${CONFIG_FILE}" <<EOF

[mcp]
enabled = true
deferred_loading = false
EOF
    fi

    if grep -q '^\[autonomy\]' "${CONFIG_FILE}" 2>/dev/null; then
        sed -i '/^\[autonomy\]/,/^\[/{s/^level = "supervised"$/level = "full"/}' "${CONFIG_FILE}"
        sed -i '/^\[autonomy\]/,/^\[/{s/^require_approval_for_medium_risk = true$/require_approval_for_medium_risk = false/}' "${CONFIG_FILE}"
    fi

    HA_TOKEN="${SUPERVISOR_TOKEN}"
    bashio::log.info "Writing Home Assistant MCP server configuration..."
    cat >> "${CONFIG_FILE}" <<EOF

[[mcp.servers]]
name = "home-assistant"
transport = "http"
url = "http://supervisor/core/api/mcp"
headers = { Authorization = "Bearer ${HA_TOKEN}" }
EOF
fi

# ── Gateway: bind to 0.0.0.0 for direct port access ──────────────────────────
GATEWAY_PAIRING=$(bashio::config 'gateway_pairing')
PAIRING_VAL="true"
if bashio::var.false "${GATEWAY_PAIRING}"; then PAIRING_VAL="false"; fi

if grep -q '^\[gateway\]' "${CONFIG_FILE}" 2>/dev/null; then
    sed -i '/^\[gateway\]/,/^\[/{s/^host = ".*"/host = "0.0.0.0"/}' "${CONFIG_FILE}"
    sed -i '/^\[gateway\]/,/^\[/{s/^port = .*/port = 42617/}' "${CONFIG_FILE}"
    sed -i '/^\[gateway\]/,/^\[/{s/^allow_public_bind = .*$/allow_public_bind = true/}' "${CONFIG_FILE}"
    sed -i '/^\[gateway\]/,/^\[/{s/^require_pairing = .*$/require_pairing = '"${PAIRING_VAL}"'/}' "${CONFIG_FILE}"
else
    cat >> "${CONFIG_FILE}" <<EOF

[gateway]
host = "0.0.0.0"
port = 42617
allow_public_bind = true
require_pairing = ${PAIRING_VAL}
EOF
fi
bashio::log.info "Gateway configured (port=42617, pairing=${PAIRING_VAL})."

# ── Browser addon auto-detect ─────────────────────────────────────────────────
BROWSER_CDP_PORT=$(bashio::config 'browser_cdp_port')
HOST_IP=""
if command -v ip > /dev/null 2>&1; then
    HOST_IP=$(ip route show default 2>/dev/null | awk '/default/ {print $3; exit}') || true
fi
if [ -z "${HOST_IP}" ]; then HOST_IP="172.30.32.1"; fi
BROWSER_CDP_URL="ws://${HOST_IP}:${BROWSER_CDP_PORT}"
BROWSER_DETECTED=false
for i in 1 2 3 4 5 6; do
    if curl -sf "http://${HOST_IP}:${BROWSER_CDP_PORT}/json/version" > /dev/null 2>&1; then
        BROWSER_DETECTED=true
        break
    fi
    if [ "${i}" -lt 6 ]; then bashio::log.info "Waiting for browser addon... (${i}/6)"; sleep 5; fi
done
if [ "${BROWSER_DETECTED}" = "true" ]; then
    bashio::log.info "Browser addon detected — configuring agent-browser CDP at ${BROWSER_CDP_URL}"
    mkdir -p "${HOME}/.agent-browser"
    printf '{"cdp":"%s"}\n' "${BROWSER_CDP_URL}" > "${HOME}/.agent-browser/config.json"
else
    bashio::log.info "No browser addon detected — browser tools unavailable. Install browserless_chrome addon to enable."
fi

# ── Start daemon ──────────────────────────────────────────────────────────────
bashio::log.info "Starting ZeroClaw daemon..."
exec /usr/local/bin/zeroclaw daemon
