# zeroclaw/rootfs/etc/s6-overlay/s6-rc.d/zeroclaw/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -e

# ── Paths ────────────────────────────────────────────────────────────────────
export HOME="/share/zeroclaw"
export SHELL="/bin/bash"
ZEROCLAW_DATA="/share/zeroclaw/.zeroclaw"
CONFIG_FILE="${ZEROCLAW_DATA}/config.toml"
mkdir -p "${ZEROCLAW_DATA}" "${HOME}/workspace"

# Redirect ZeroClaw's default internal data dir to the persistent share volume.
# v0.6.8+ uses /usr/local/var/zeroclaw/ instead of ~/.zeroclaw/ — symlink it.
mkdir -p /usr/local/var
ln -sfn "${ZEROCLAW_DATA}" /usr/local/var/zeroclaw

# ── Read options ──────────────────────────────────────────────────────────────
API_KEY=$(bashio::config 'api_key')
TIMEZONE=$(bashio::config 'timezone')

# ── Validate required options ─────────────────────────────────────────────────
if bashio::var.is_empty "${API_KEY}"; then
    bashio::log.fatal "API key is required. Set it in the add-on configuration."
    exit 1
fi

# ── First-run: bootstrap config via zeroclaw onboard ─────────────────────────
if [ ! -f "${CONFIG_FILE}" ]; then
    bashio::log.info "First start detected — running initial setup..."
    /usr/local/bin/zeroclaw onboard \
        --api-key "${API_KEY}" \
        --memory markdown \
        --force
    bashio::log.info "Setup complete. Configure provider, model and channels via the web dashboard."
fi

# ── Export API key for daemon and terminal sessions ───────────────────────────
export ZEROCLAW_API_KEY="${API_KEY}"
cat > /etc/profile.d/zeroclaw-api.sh <<EOF
export ZEROCLAW_API_KEY="${API_KEY}"
EOF

# ── Timezone ──────────────────────────────────────────────────────────────────
if [ -f "/usr/share/zoneinfo/${TIMEZONE}" ]; then
    ln -sf "/usr/share/zoneinfo/${TIMEZONE}" /etc/localtime
    echo "${TIMEZONE}" > /etc/timezone
    bashio::log.info "Timezone set to ${TIMEZONE}"
else
    bashio::log.warning "Unknown timezone '${TIMEZONE}', falling back to UTC."
fi

# ── Remove managed MCP section before re-injection ───────────────────────────
# [[mcp.servers]] is always appended — truncate at the marker so it is
# rewritten on every start, keeping ha_mcp_enabled as source of truth.
# NOTE: grep exits 1 on no match — use || true to avoid aborting with set -e.
line=$(grep -n '^\[\[mcp\.servers\]\]' "${CONFIG_FILE}" 2>/dev/null | head -1 | cut -d: -f1 || true)
if [ -n "${line}" ] && [ "${line}" -gt 1 ]; then
    head -n $((line - 1)) "${CONFIG_FILE}" > "${CONFIG_FILE}.tmp"
    mv "${CONFIG_FILE}.tmp" "${CONFIG_FILE}"
fi

# ── Home Assistant MCP integration ────────────────────────────────────────────
HA_MCP_ENABLED=$(bashio::config 'ha_mcp_enabled')
if bashio::var.true "${HA_MCP_ENABLED}"; then
    if grep -q '^\[mcp\]' "${CONFIG_FILE}" 2>/dev/null; then
        sed -i '/^\[mcp\]/,/^\[/{s/^enabled = false$/enabled = true/}' "${CONFIG_FILE}"
        sed -i '/^\[mcp\]/,/^\[/{s/^deferred_loading = true$/deferred_loading = false/}' "${CONFIG_FILE}"
        sed -i '/^servers = \[\]$/d' "${CONFIG_FILE}"
        bashio::log.info "MCP section updated: enabled = true, deferred_loading = false."
    else
        cat >> "${CONFIG_FILE}" <<EOF

[mcp]
enabled = true
deferred_loading = false
EOF
    fi

    if grep -q '^\[autonomy\]' "${CONFIG_FILE}" 2>/dev/null; then
        sed -i '/^\[autonomy\]/,/^\[/{s/^level = "supervised"$/level = "full"/}' "${CONFIG_FILE}"
        sed -i '/^\[autonomy\]/,/^\[/{s/^require_approval_for_medium_risk = true$/require_approval_for_medium_risk = false/}' "${CONFIG_FILE}"
    fi

    HA_TOKEN="${SUPERVISOR_TOKEN}"
    bashio::log.info "Writing Home Assistant MCP server configuration..."
    cat >> "${CONFIG_FILE}" <<EOF

[[mcp.servers]]
name = "home-assistant"
transport = "http"
url = "http://supervisor/core/api/mcp"

[mcp.servers.headers]
Authorization = "Bearer ${HA_TOKEN}"
EOF
fi

# ── Gateway: bind to 0.0.0.0 for direct port access ──────────────────────────
GATEWAY_PAIRING=$(bashio::config 'gateway_pairing')
PAIRING_VAL="true"
if bashio::var.false "${GATEWAY_PAIRING}"; then PAIRING_VAL="false"; fi

if grep -q '^\[gateway\]' "${CONFIG_FILE}" 2>/dev/null; then
    sed -i '/^\[gateway\]/,/^\[/{s/^host = ".*"/host = "0.0.0.0"/}' "${CONFIG_FILE}"
    sed -i '/^\[gateway\]/,/^\[/{s/^port = .*/port = 42617/}' "${CONFIG_FILE}"
    sed -i '/^\[gateway\]/,/^\[/{s/^allow_public_bind = .*$/allow_public_bind = true/}' "${CONFIG_FILE}"
    sed -i '/^\[gateway\]/,/^\[/{s/^require_pairing = .*$/require_pairing = '"${PAIRING_VAL}"'/}' "${CONFIG_FILE}"
else
    cat >> "${CONFIG_FILE}" <<EOF

[gateway]
host = "0.0.0.0"
port = 42617
allow_public_bind = true
require_pairing = ${PAIRING_VAL}
EOF
fi
bashio::log.info "Gateway configured (port=42617, pairing=${PAIRING_VAL})."

# ── Start daemon ──────────────────────────────────────────────────────────────
bashio::log.info "Starting ZeroClaw daemon..."
exec /usr/local/bin/zeroclaw daemon
