# zeroclaw/rootfs/etc/s6-overlay/s6-rc.d/zeroclaw/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -e

# ── Paths ────────────────────────────────────────────────────────────────────
export HOME="/share/zeroclaw"
export SHELL="/bin/bash"

ZEROCLAW_DATA="/share/zeroclaw/.zeroclaw"
CONFIG_FILE="${ZEROCLAW_DATA}/config.toml"
# v0.8.0+ resolves the install root from ZEROCLAW_CONFIG_DIR (priority 1) —
# explicit beats relying on HOME/legacy fallbacks. Pre-create the V3 layout:
# data/ (shared DBs), shared/skills/ (host-wide skills), and the per-agent
# workspace for the implicit "default" agent created by V2→V3 migration.
export ZEROCLAW_CONFIG_DIR="${ZEROCLAW_DATA}"
mkdir -p \
    "${ZEROCLAW_DATA}/data" \
    "${ZEROCLAW_DATA}/shared/skills" \
    "${ZEROCLAW_DATA}/agents/default/workspace"

# ── Inject environment variables from HA UI configuration ────────────────────
# Source via jq from /data/options.json (single source of truth, deduped).
if [ -f /data/options.json ]; then
    while IFS='=' read -r _name _value; do
        [ -z "${_name}" ] && continue
        export "${_name}=${_value}"
        printf '%s' "${_value}" > "/var/run/s6/container_environment/${_name}" 2>/dev/null || true
        bashio::log.info "Env var set: ${_name}"
    done < <(jq -r '(.env_vars // []) | map(select(.name != null and .name != "")) | unique_by(.name) | .[] | "\(.name)=\(.value)"' /data/options.json)
fi

# ── Write env vars to workspace so agents can read them via file_read ─────────
# (shell tool cannot access env vars under workspace_only security policy)
# Build via jq from /data/options.json — single source of truth, no bashio loop
# (bashio::config iteration was observed duplicating entries 32x at runtime).
# V3 splits the workspace per agent: write under the default agent's workspace,
# which is where the V2→V3 filesystem migration moves the legacy workspace/.
WORKSPACE_CONFIG="${ZEROCLAW_DATA}/agents/default/workspace/config"
mkdir -p "${WORKSPACE_CONFIG}"
ENV_VARS_FILE="${WORKSPACE_CONFIG}/env_vars.conf"
if [ -f /data/options.json ]; then
    jq -r '(.env_vars // []) | map(select(.name != null and .name != "")) | unique_by(.name) | .[] | "\(.name)=\(.value)"' \
        /data/options.json > "${ENV_VARS_FILE}.tmp" && \
        mv "${ENV_VARS_FILE}.tmp" "${ENV_VARS_FILE}"
else
    : > "${ENV_VARS_FILE}"
fi
chmod 600 "${ENV_VARS_FILE}"

# ── First-run guard ───────────────────────────────────────────────────────────
if [ ! -f "${CONFIG_FILE}" ]; then
    bashio::log.fatal "No configuration found. Open the web terminal and run: zeroclaw onboard"
    exit 1
fi

# ── V3-only ──────────────────────────────────────────────────────────────────
# This add-on supports schema_version = 3 only. Upgrades from V2 (ZeroClaw
# 0.7.x) are not handled here: run `zeroclaw config migrate` from the web
# terminal once before installing 0.8.x. We don't trust the upstream migrator
# enough to wrap it — see issues #16/#18 for why.
#
# Boot pipeline:
#   A — read HA add-on options
#   B — V3 sanitize: drop the managed [[mcp.servers]] entry so Phase D can
#       re-inject it idempotently with a fresh Supervisor token
#   C — Gateway settings (always)
#   D — Home Assistant MCP inject (only if ha_mcp_enabled = true)
#   E — Browser CDP auto-detect (always)
#   F — Validation via `zeroclaw config list` — fail fast on bad TOML
#   G — exec the daemon

# ── Phase A: read HA add-on options ───────────────────────────────────────────
TIMEZONE=$(bashio::config 'timezone')
HA_MCP_ENABLED=$(bashio::config 'ha_mcp_enabled')
GATEWAY_PAIRING=$(bashio::config 'gateway_pairing')
BROWSER_CDP_PORT=$(bashio::config 'browser_cdp_port')
BROWSER_STEALTH=$(bashio::config 'browser_stealth')

if [ -f "/usr/share/zoneinfo/${TIMEZONE}" ]; then
    ln -sf "/usr/share/zoneinfo/${TIMEZONE}" /etc/localtime
    echo "${TIMEZONE}" > /etc/timezone
    bashio::log.info "Timezone set to ${TIMEZONE}"
else
    bashio::log.warning "Unknown timezone '${TIMEZONE}', falling back to UTC."
fi

# ── Phase B: V3 sanitize ──────────────────────────────────────────────────────
# Drop the previous boot's managed home-assistant entry so Phase D can re-inject
# it cleanly with the current Supervisor token. AOT (`[[mcp.servers]]`) is the
# only shape the V3 deserializer reads from disk, so that's the only one we
# write and the only one we strip. User-defined entries under any other name
# are preserved.
MCP_HA="${ZEROCLAW_DATA}/mcp-ha.toml"
awk '
function flush_block() {
    if (block_buf != "" && !drop_block) printf "%s", block_buf
    block_buf = ""; drop_block = 0
}
BEGIN { block_buf = ""; drop_block = 0; in_aot = 0 }
/^\[/ {
    flush_block()
    block_buf = $0 "\n"
    in_aot = ($0 == "[[mcp.servers]]")
    next
}
in_aot && /^name = "home-assistant"$/ { drop_block = 1 }
{ block_buf = block_buf $0 "\n" }
END { flush_block() }
' "${CONFIG_FILE}" > "${CONFIG_FILE}.tmp" && mv "${CONFIG_FILE}.tmp" "${CONFIG_FILE}"

# ── Phase C: Gateway settings ─────────────────────────────────────────────────
PAIRING_VAL="true"
if bashio::var.false "${GATEWAY_PAIRING}"; then PAIRING_VAL="false"; fi

WEB_DIST_DIR="/opt/zeroclaw/web/dist"
if grep -q '^\[gateway\]' "${CONFIG_FILE}" 2>/dev/null; then
    sed -i '/^\[gateway\]/,/^\[/{s/^host = ".*"/host = "0.0.0.0"/}' "${CONFIG_FILE}"
    sed -i '/^\[gateway\]/,/^\[/{s/^port = .*/port = 42617/}' "${CONFIG_FILE}"
    sed -i '/^\[gateway\]/,/^\[/{s/^allow_public_bind = .*$/allow_public_bind = true/}' "${CONFIG_FILE}"
    sed -i '/^\[gateway\]/,/^\[/{s/^require_pairing = .*$/require_pairing = '"${PAIRING_VAL}"'/}' "${CONFIG_FILE}"
    if awk '/^\[gateway\]/{f=1;next} /^\[/{f=0} f' "${CONFIG_FILE}" | grep -q '^web_dist_dir = '; then
        sed -i '/^\[gateway\]/,/^\[/{s|^web_dist_dir = .*|web_dist_dir = "'"${WEB_DIST_DIR}"'"|}' "${CONFIG_FILE}"
    else
        sed -i '/^\[gateway\]/a web_dist_dir = "'"${WEB_DIST_DIR}"'"' "${CONFIG_FILE}"
    fi
else
    cat >> "${CONFIG_FILE}" <<EOF

[gateway]
host = "0.0.0.0"
port = 42617
allow_public_bind = true
require_pairing = ${PAIRING_VAL}
web_dist_dir = "${WEB_DIST_DIR}"
EOF
fi
bashio::log.info "Gateway configured (port=42617, pairing=${PAIRING_VAL})."

# ── Phase D: Home Assistant MCP inject ────────────────────────────────────────
# V3 defaults `[mcp]` to enabled=true / deferred_loading=false, so we only
# append the section if the user hasn't already declared it. Then attach the
# managed server entry. AOT form: the V3 deserializer only loads
# `[[mcp.servers]]`, never dotted.
if bashio::var.true "${HA_MCP_ENABLED}"; then
    if ! grep -q '^\[mcp\]' "${CONFIG_FILE}"; then
        cat >> "${CONFIG_FILE}" <<EOF

[mcp]
enabled = true
deferred_loading = false
EOF
    fi

    HA_TOKEN="${SUPERVISOR_TOKEN}"
    cat > "${MCP_HA}" <<EOF

[[mcp.servers]]
name = "home-assistant"
transport = "http"
url = "http://supervisor/core/api/mcp"
headers = { Authorization = "Bearer ${HA_TOKEN}" }
EOF
    cat "${MCP_HA}" >> "${CONFIG_FILE}"
    bashio::log.info "Home Assistant MCP server configured."
else
    rm -f "${MCP_HA}"
fi

# ── Phase E: Browser addon auto-detect ────────────────────────────────────────
HOST_IP=""
if command -v ip > /dev/null 2>&1; then
    HOST_IP=$(ip route show default 2>/dev/null | awk '/default/ {print $3; exit}') || true
fi
if [ -z "${HOST_IP}" ]; then HOST_IP="172.30.32.1"; fi
BROWSER_CDP_URL="ws://${HOST_IP}:${BROWSER_CDP_PORT}"
if bashio::var.true "${BROWSER_STEALTH}"; then
    BROWSER_CDP_URL="${BROWSER_CDP_URL}?stealth=true"
fi
BROWSER_DETECTED=false
for i in 1 2 3 4 5 6; do
    if curl -sf -o /dev/null --max-time 3 "http://${HOST_IP}:${BROWSER_CDP_PORT}/sessions" \
        || curl -sf -o /dev/null --max-time 3 "http://${HOST_IP}:${BROWSER_CDP_PORT}/json/version"; then
        BROWSER_DETECTED=true
        break
    fi
    if [ "${i}" -lt 6 ]; then bashio::log.info "Waiting for browser addon... (${i}/6)"; sleep 5; fi
done
if [ "${BROWSER_DETECTED}" = "true" ]; then
    bashio::log.info "Browser addon detected — configuring agent-browser CDP at ${BROWSER_CDP_URL}"
    mkdir -p "${HOME}/.agent-browser"
    printf '{"cdp":"%s"}\n' "${BROWSER_CDP_URL}" > "${HOME}/.agent-browser/config.json"
else
    bashio::log.info "No browser addon detected — browser tools unavailable. Install browserless_chrome addon to enable."
fi

# ── Phase F: Validation ───────────────────────────────────────────────────────
# `zeroclaw config list` loads, parses, and deserializes the entire schema.
# If the file is broken (orphan stripping went wrong, an inject produced an
# invalid shape, a hand-edit left dangling syntax), this command exits non-zero
# and we fail loudly with a pointer to the file. Without this gate the daemon
# would refuse to serve with the cryptic "malformed security-critical sections
# (<entire-config>)" message and we'd see a crash-loop instead of a clear
# error. Backup is dropped beside the live file so the operator can restore.
if ! /usr/local/bin/zeroclaw --config-dir "${ZEROCLAW_DATA}" config list >/dev/null 2>&1; then
    cp -a "${CONFIG_FILE}" "${CONFIG_FILE}.invalid-$(date -u +%Y%m%dT%H%M%SZ)" || true
    bashio::log.fatal "Config validation failed after sanitize/inject. \
A copy of the invalid file was saved next to ${CONFIG_FILE} (suffix '.invalid-*'). \
Run \`zeroclaw config migrate\` from the web terminal to inspect the precise error, \
then restart the add-on."
    exit 1
fi
bashio::log.info "Config validation passed."

# ── Phase G: Start daemon ─────────────────────────────────────────────────────
# Pass --config-dir explicitly. The `export ZEROCLAW_CONFIG_DIR=...` above
# does not survive `exec` on this s6-overlay/Docker combination (HOME gets
# reset to /root and ZEROCLAW_CONFIG_DIR is dropped), which makes the daemon
# fall back to `${HOME}/.zeroclaw/` for `install_root_dir()` — gateway
# personality endpoints then look up SOUL.md/AGENTS.md in container-ephemeral
# paths instead of /share/zeroclaw/.zeroclaw/agents/<alias>/workspace/.
# `--config-dir` is a global CLI flag honored before any config resolution
# runs, so it's the single source of truth (see issue #22).
bashio::log.info "Starting ZeroClaw daemon..."
exec /usr/local/bin/zeroclaw --config-dir "${ZEROCLAW_DATA}" daemon
