# zeroclaw/rootfs/etc/s6-overlay/s6-rc.d/zeroclaw/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -e

# ── Paths ────────────────────────────────────────────────────────────────────
export HOME="/share/zeroclaw"
CONFIG_FILE="${HOME}/.zeroclaw/config.toml"
mkdir -p "${HOME}/.zeroclaw" "${HOME}/workspace"

# ── Read options ──────────────────────────────────────────────────────────────
PROVIDER=$(bashio::config 'provider')
API_KEY=$(bashio::config 'api_key')
MODEL=$(bashio::config 'model')
TIMEZONE=$(bashio::config 'timezone')

# ── Validate API key (with fallback to stored config on reinstall) ────────────
if bashio::var.is_empty "${API_KEY}"; then
    if [ -f "${CONFIG_FILE}" ] && grep -q '^api_key = ' "${CONFIG_FILE}"; then
        bashio::log.info "API key not set in UI — using stored key from config.toml (reinstall detected)."
        API_KEY=""  # zeroclaw reads it from config.toml directly
    else
        bashio::log.fatal "API key is required. Set it in the add-on configuration."
        exit 1
    fi
fi

# ── First-run: bootstrap config via zeroclaw onboard ─────────────────────────
if [ ! -f "${CONFIG_FILE}" ]; then
    bashio::log.info "First start detected — running initial setup..."

    PROVIDER_ARG="${PROVIDER}"
    # 'other' is a UI placeholder; let onboard use the key without --provider
    if [ "${PROVIDER_ARG}" = "other" ]; then
        PROVIDER_ARG=""
    fi

    if [ -n "${PROVIDER_ARG}" ]; then
        /usr/local/bin/zeroclaw onboard \
            --provider "${PROVIDER_ARG}" \
            --api-key "${API_KEY}" \
            --model "${MODEL}" \
            --memory markdown \
            --force
    else
        /usr/local/bin/zeroclaw onboard \
            --api-key "${API_KEY}" \
            --memory markdown \
            --force
    fi

    bashio::log.info "Setup complete."
fi

# ── Export API key for daemon and write to profile for terminal sessions ──────
if [ -n "${API_KEY}" ]; then
    export ZEROCLAW_API_KEY="${API_KEY}"

    # Map provider to its specific env var so 'zeroclaw agent' in the terminal works
    case "${PROVIDER}" in
        gemini|google)    PROVIDER_KEY_VAR="GEMINI_API_KEY" ;;
        openai)           PROVIDER_KEY_VAR="OPENAI_API_KEY" ;;
        anthropic)        PROVIDER_KEY_VAR="ANTHROPIC_API_KEY" ;;
        openrouter)       PROVIDER_KEY_VAR="OPENROUTER_API_KEY" ;;
        mistral)          PROVIDER_KEY_VAR="MISTRAL_API_KEY" ;;
        groq)             PROVIDER_KEY_VAR="GROQ_API_KEY" ;;
        deepseek)         PROVIDER_KEY_VAR="DEEPSEEK_API_KEY" ;;
        xai)              PROVIDER_KEY_VAR="XAI_API_KEY" ;;
        *)                PROVIDER_KEY_VAR="ZEROCLAW_API_KEY" ;;
    esac
    export "${PROVIDER_KEY_VAR}=${API_KEY}"

    # Write to profile so bash login shells (ttyd terminal) inherit the key
    cat > /etc/profile.d/zeroclaw-api.sh <<EOF
export ZEROCLAW_API_KEY="${API_KEY}"
export ${PROVIDER_KEY_VAR}="${API_KEY}"
EOF
fi

# ── Timezone ──────────────────────────────────────────────────────────────────
if [ -f "/usr/share/zoneinfo/${TIMEZONE}" ]; then
    ln -sf "/usr/share/zoneinfo/${TIMEZONE}" /etc/localtime
    echo "${TIMEZONE}" > /etc/timezone
    bashio::log.info "Timezone set to ${TIMEZONE}"
else
    bashio::log.warning "Unknown timezone '${TIMEZONE}', falling back to UTC."
fi

# ── Remove managed sections before re-injection ───────────────────────────────
# Telegram and MCP server entries are always appended at the end of config.toml.
# Truncate at the first managed section marker so they are rewritten on every
# start, keeping HA UI options as the single source of truth.
# NOTE: grep exits 1 when no match — use || true to avoid aborting with set -e.
for marker in '^\[channels_config\.telegram\]' '^\[\[mcp\.servers\]\]'; do
    line=$(grep -n "${marker}" "${CONFIG_FILE}" 2>/dev/null | head -1 | cut -d: -f1 || true)
    if [ -n "${line}" ] && [ "${line}" -gt 1 ]; then
        head -n $((line - 1)) "${CONFIG_FILE}" > "${CONFIG_FILE}.tmp"
        mv "${CONFIG_FILE}.tmp" "${CONFIG_FILE}"
    fi
done

# ── Telegram integration ───────────────────────────────────────────────────────
TOKEN=$(bashio::config 'telegram_token')
ALLOWED=$(bashio::config 'telegram_allowed_users')
if [ -n "${TOKEN}" ]; then
    # Build allowed_users array: "*" becomes ["*"], else split CSV into array
    if [ "${ALLOWED}" = "*" ]; then
        ALLOWED_TOML='["*"]'
    else
        ALLOWED_TOML="[$(echo "${ALLOWED}" | sed 's/,/","/g; s/^/"/; s/$/"/' )]"
    fi

    bashio::log.info "Writing Telegram channel configuration..."
    cat >> "${CONFIG_FILE}" <<EOF

[channels_config.telegram]
bot_token = "${TOKEN}"
allowed_users = ${ALLOWED_TOML}
EOF
fi

# ── Home Assistant MCP integration ────────────────────────────────────────────
HA_MCP_ENABLED=$(bashio::config 'ha_mcp_enabled')
if bashio::var.true "${HA_MCP_ENABLED}"; then
    # Ensure [mcp] section exists with enabled = true
    if grep -q '^\[mcp\]' "${CONFIG_FILE}" 2>/dev/null; then
        # onboard sets enabled = false and servers = [] — fix both
        sed -i '/^\[mcp\]/,/^\[/{s/^enabled = false$/enabled = true/}' "${CONFIG_FILE}"
        sed -i '/^servers = \[\]$/d' "${CONFIG_FILE}"
        bashio::log.info "MCP section updated: enabled = true."
    else
        cat >> "${CONFIG_FILE}" <<EOF

[mcp]
enabled = true
deferred_loading = true
EOF
    fi

    # Use HA-provided supervisor token — no manual token required
    HA_TOKEN="${SUPERVISOR_TOKEN}"
    bashio::log.info "Writing Home Assistant MCP server configuration..."
    cat >> "${CONFIG_FILE}" <<EOF

[[mcp.servers]]
name = "home-assistant"
transport = "http"
url = "http://supervisor/core/api/mcp"

[mcp.servers.headers]
Authorization = "Bearer ${HA_TOKEN}"
EOF
fi

# ── Start daemon ──────────────────────────────────────────────────────────────
bashio::log.info "Starting ZeroClaw daemon (provider=${PROVIDER}, model=${MODEL})..."
exec /usr/local/bin/zeroclaw daemon
