# nullclaw/rootfs/etc/s6-overlay/s6-rc.d/nullclaw/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "NullClaw gateway stopped."
