# zeroclaw/rootfs/etc/s6-overlay/s6-rc.d/lightpanda/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "Lightpanda headless browser stopped."
