# zeroclaw/rootfs/etc/s6-overlay/s6-rc.d/lightpanda/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -e
bashio::log.info "Starting Lightpanda headless browser (CDP) on port 9222..."
exec /usr/local/bin/lightpanda serve --host 127.0.0.1 --port 9222 --timeout 604800
