# picoclaw/rootfs/etc/s6-overlay/s6-rc.d/picoclaw/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "PicoClaw gateway stopped."
